# cookbooks/beater/templates/default/network-wireless@.service.erb
[Unit]
Description=Wireless network connectivity (%i)
Wants=network.target
Before=network.target
BindsTo=sys-subsystem-net-devices-%i.device
After=sys-subsystem-net-devices-%i.device

[Service]
Type=oneshot
RemainAfterExit=yes

ExecStart=/usr/sbin/ip link set dev %i up
ExecStop=/usr/sbin/ip link set dev %i down

[Install]
WantedBy=multi-user.target
